# indexer.cron
#!/bin/bash
#
# description: musicbox indexer cronjob
#

# Change this to your install directory
export MUSICBOX_INSTALL_DIR=/opt/musicbox
export MYSQL_HOST=db1.rosa.com
export MYSQL_USER=musicbox
export MYSQL_PASS=rosa
export MYSQL_DB=musicbox

"$MUSICBOX_INSTALL_DIR/indexer.exp"

